# roles/podman_rootless/templates/wp-stack-redis.container.j2
[Unit]
Description=Redis Container
After=wp-stack.pod

[Container]
Image=docker.io/library/redis:alpine
# Add Pod
Pod=wp-stack.pod

[Service]
Restart=always
